# etc/profile.d/site.sh
CONFIG_SITE="/usr/share/site/x86_64-unknown-linux-gnu"
export CONFIG_SITE
