# etc/udev/rules.d/99-iwlwifi-led.rules
SUBSYSTEM=="leds", ACTION=="add", KERNEL=="iwl-phy*:assoc", RUN+="/lib/udev/iwlwifi-led.sh"

